# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:backend]
command=python /app/backend/app.py
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
environment=PYTHONUNBUFFERED=1
